# deploy/systemd/quantpipe-pipeline.timer
[Unit]
Description=QuantPipe Daily Pipeline — Mon-Fri 21:30 UTC (post-US-market-close)
Requires=quantpipe-pipeline.service

[Timer]
# 21:30 UTC = ~30 min after US market close (4:00 PM ET summer / 5:00 PM ET winter).
# Running post-close ensures yfinance returns same-day closing prices rather than
# yesterday's close, so the dashboard always shows fresh data by evening.
OnCalendar=Mon-Fri 21:30:00 UTC

# Do NOT catch up missed runs on daemon-reload / reboot.
# Persistent=true caused update.sh's daemon-reload to fire a surprise
# mid-day pipeline run; keeping it false means the next scheduled window
# is used instead.
Persistent=false

# Tiny jitter (0-60 s) so retries don't pile up on reboot
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
