# deploy/systemd/quantpipe-rebalance.service
[Unit]
Description=QuantPipe IBKR Paper Rebalance (post-pipeline execution)
After=network-online.target quantpipe-pipeline.service
Wants=network-online.target

[Service]
Type=oneshot
User=quantpipe
WorkingDirectory=/opt/quantpipe
EnvironmentFile=/opt/quantpipe/.env
Environment="PATH=/opt/quantpipe/.venv/bin:/home/quantpipe/.local/bin:/usr/local/bin:/usr/bin:/bin"

ExecStart=/opt/quantpipe/.venv/bin/python orchestration/rebalance.py \
    --broker ibkr

# Exit code 1 (pre-trade blocked) is still a clean exit — don't treat as failure
SuccessExitStatus=0 1
Restart=no

StandardOutput=append:/var/log/quantpipe/rebalance.log
StandardError=append:/var/log/quantpipe/rebalance.log

MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
